# ops/systemd/xcore-metrics-gateway.service
[Unit]
Description=XCore Metrics Gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=xcore
Group=xcore
WorkingDirectory=/opt/xcore/xcore-metrics-gateway
EnvironmentFile=/etc/xcore/xcore-metrics-gateway.env
ExecStart=/usr/local/bin/uv run --project /opt/xcore/xcore-metrics-gateway xcore-metrics-gateway
Restart=on-failure
RestartSec=5
TimeoutStopSec=15
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
